# conf/wol-proxy.service.template
[Unit]
Description=WOL Proxy

[Service]
User=alex
Environment="WOL_LOG=/opt/wol-proxy/log"
ExecStart=/usr/bin/java -jar /opt/wol-proxy/wol.jar --server @SERVER_PORT@ @PROXY_PORT@
SuccessExitStatus=143
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
